# roles/compositional/files/portal_commands_receivable.service
[Service]
StandardOutput=journal+console
StandardError=journal+console
ExecStart=/usr/local/bin/commands_receivable.py
