# roles/compositional/files/portal_commands_receivable.socket
[Socket]
ListenStream=/var/run/commands_receiveable.sock

[Install]
WantedBy=sockets.target
