# roles/compositional/files/portal_commands_recievable.socket
[Socket]
ListenStream=/var/run/commands_recievable.sock

[Install]
WantedBy=sockets.target
