# roles/compositional/files/portal_commands_receivable.socket
[Socket]
ListenStream=/var/run/commands_receivable.sock

[Install]
WantedBy=sockets.target
