# roles/compositional/files/portal_commands_recievable.service
[Service]
StandardOutput=journal+console
StandardError=journal+console
ExecStart=/root/commands_recievable.py
